# dnspodsh/dnspod.service
[Unit]
Description=DDNS

[Service]
Type=simple
ExecStart=/root/Script/dnspodsh/dnspodsh.sh &>/dev/null
