# acme.sh/acme.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/root/Script/acme.sh &>/dev/null